# files/etc/systemd/system/netvision-kiosk.service
[Unit]
Description=NetVisionConnect Kiosk (Chromium plein écran)
After=systemd-user-sessions.service network-online.target
Wants=network-online.target

[Service]
User=kiosk
PAMName=login
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=journal
StandardError=journal
# Lance un serveur X minimal qui exécute le script kiosque
ExecStart=/usr/bin/xinit /usr/local/bin/netvision-kiosk.sh -- :0 vt1 -nocursor -keeptty
# Watchdog logiciel : si la session tombe, on la relance automatiquement
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
