# files/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
# Connexion automatique de l'utilisateur "kiosk" sur la console tty1
ExecStart=
ExecStart=-/sbin/agetty --autologin kiosk --noclear %I $TERM
